# deploy/onprem/scheduler/crontab
# KlinoVax on-prem cron programı (supercronic formatı).
# Vercel Cron'un birebir karşılığı — her satır app'in /api/cron/* rotasını
# CRON_SECRET Bearer'ı ile tetikler. Saatler UTC.
# Altyapı cron'ları (cleanup/backup/verify-backup/license-heartbeat) her durumda
# çalışır; iş cron'ları READONLY/LOCKED'ta route içinde kendini atlar.

# ── Altyapı ──
0 3 * * *    curl -fsS -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/cleanup
15 3 * * *   curl -fsS -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/backup
0 4 * * *    curl -fsS -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/verify-backup
0 */6 * * *  curl -fsS -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/license-heartbeat

# ── İş ──
30 5 * * *   curl -fsS -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/staff-sync
0 7 * * *    curl -fsS -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/reminders
15 4 * * *   curl -fsS -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/expo-receipts
0 8 * * *    curl -fsS -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/subscription-reminders
0 6 * * *    curl -fsS -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/daily-quiz-push
5 21 * * *   curl -fsS -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/streak-maintenance
0 8 1 * *    curl -fsS -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/compliance-report
0 1 1 1 *    curl -fsS -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/period-rollover
