# deploy/onprem/scheduler/crontab
# KlinoVax on-prem cron programı (supercronic formatı).
# Vercel Cron'un birebir karşılığı — her satır app'in /api/cron/* rotasını
# CRON_SECRET Bearer'ı ile tetikler. Saatler UTC (TR = UTC+3; ör. 03:00 UTC = 06:00 TRT).
# Altyapı cron'ları (cleanup/backup/verify-backup/license-heartbeat) her durumda
# çalışır; iş cron'ları READONLY/LOCKED'ta route içinde kendini atlar.
#
# RETRY: app o dakika ayakta değilse (güncelleme penceresi/OOM-restart/geçici hata) iş
# SESSİZCE atlanmasın diye her curl retry'lı. --retry-all-errors = bağlantı-reddi + 5xx
# dahil TÜM hatalarda tekrar dener (curl 7.71+; bookworm curl 7.88 destekler). 5×90s ≈ 7.5 dk
# tolerans penceresi. Yedek/heartbeat gibi altyapı işlerinin app-restart'a denk gelip
# atlanması (44-günlük sessiz yedek olayının kardeşi) bu şekilde önlenir.

RETRY=--retry 5 --retry-delay 90 --retry-all-errors --max-time 600

# ── Liveness beacon ──
# Supercronic'in canlı VE tetikliyor olduğunu kanıtlar: her 5 dk /tmp/scheduler-alive'a
# epoch yazar. Compose healthcheck bu dosyanın tazeliğine bakar (bkz. docker-compose.yml
# scheduler.healthcheck) — bayatlarsa container 'unhealthy'. supercronic komutu /bin/sh -c
# ile çalıştırır → redirection çalışır. Ağ/app gerektirmez; saf yerel yazım.
*/5 * * * *  date +%s > /tmp/scheduler-alive

# ── Altyapı ──
0 3 * * *    curl -fsS $RETRY -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/cleanup
15 3 * * *   curl -fsS $RETRY -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/backup
0 4 * * *    curl -fsS $RETRY -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/verify-backup
0 */6 * * *  curl -fsS $RETRY -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/license-heartbeat

# ── İş ──
30 5 * * *   curl -fsS $RETRY -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/staff-sync
0 7 * * *    curl -fsS $RETRY -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/reminders
15 4 * * *   curl -fsS $RETRY -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/expo-receipts
0 8 * * *    curl -fsS $RETRY -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/subscription-reminders
0 6 * * *    curl -fsS $RETRY -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/daily-quiz-push
5 21 * * *   curl -fsS $RETRY -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/streak-maintenance
0 8 1 * *    curl -fsS $RETRY -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/compliance-report
0 1 1 1 *    curl -fsS $RETRY -H "Authorization: Bearer $CRON_SECRET" $APP_URL/api/cron/period-rollover
